# supervisord.conf
[program:webserver]
command=python /app/webserver_pi/tornado_pi/tornado_app.py


[program:humidity]
command=python /app/humidity_pi/humidity_logger.py